# rc.local
#!/bin/sh -e

iptables -F INPUT &
exit 0
