# rc.local
#!/bin/sh -e

su - ubuntu -c "~/TF2C_Server/server.sh start" &
exit 0
